# node_agent/deploy/llm-gateway-agent.service
[Unit]
Description=LLM Gateway — Node Agent (DGX Spark)
After=network.target
Wants=network.target

[Service]
Type=simple
User=llmservice
Group=llmservice

WorkingDirectory=/opt/llm-gateway/node_agent
EnvironmentFile=/etc/llm-gateway-agent/env
ExecStart=/opt/llm-gateway/venv-agent/bin/uvicorn main:app \
    --host 0.0.0.0 \
    --port 9443 \
    --workers 1 \
    --loop uvloop \
    --ssl-keyfile /etc/llm-gateway-agent/tls/agent.key \
    --ssl-certfile /etc/llm-gateway-agent/tls/agent.crt

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=llm-gateway-agent

# Sécurité système
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/llm-gateway-agent /tmp
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
